# resources/ansible/roles/pwiki/templates/iptables-autoexpiry.timer.j2
[Unit]
Description=Run iptables auto-expiry hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
