# crontab/root.crontab
# m h  dom mon dow   command
0   3  *   *   4     cronic bash -c '/home/pwiki/pwiki/scripts/backup-database.sh /etc/pwiki/pwiki-secrets "/home/pwikibackup/backups/pwiki-$(date '+%Y-%m-%d').sql.xz"'
0   3  *   *   5     cronic find /home/pwikibackup/backups -type f -mtime +30 -delete
